# src/dbxtool.service
[Unit]
Description=Secure Boot DBX updater

[Service]
Type=oneshot
ExecStart=/usr/bin/dbxtool -a /usr/share/dbxtool/*.bin
